# base-env/images/ubuntu-devtools/root/.sprite/etc/profile.d/languages.sh
#!/bin/bash
# PATH setup for language-specific binaries and global packages

# npm global bin directory is appended below by nodejs.sh
