# base-env/images/ubuntu-devtools/root/.sprite/etc/profile.d/sprite.sh
# Add /.sprite/bin to PATH if not already present
case ":$PATH:" in
  *":/.sprite/bin:"*) ;;
  *) export PATH="/.sprite/bin:$PATH" ;;
esac

export BROWSER="/usr/local/bin/sprite-browser"