# base-env/images/ubuntu-devtools/root/etc/profile.d/sprite.sh
export BROWSER="/usr/local/bin/sprite-browser"